# deployment/kickstart-kafka/supervisor/conf/kafka-broker.ini
[program:kafka-broker]
command=chmod 777 "/opt/kafka/script/start.sh"; /bin/bash -c "/opt/kafka/script/start.sh"
process_name=%(program_name)s
pidfile=/var/run/kafka-broker.pid
numprocs=1
autorestart=false
autostart=true
startsecs=120
startretries=3
stopsignal=TERM
stopwaitsecs=10
